# deploy/amvera-all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:xray]
command=/usr/local/bin/xray run -config /etc/xray/config.json
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=/usr/local/bin/node --enable-source-maps /app/server/index.mjs
directory=/app/server
autostart=true
autorestart=true
startretries=5
; Run the Node process as an unprivileged user. Xray stays root (port 443
; and TUN/raw-socket capabilities). supervisord itself needs root to manage
; child processes and drop privileges via per-program user= directives.
user=node
environment=NODE_ENV="production",PORT="%(ENV_PORT)s",STATIC_DIR="/app/public",XRAY_CONFIG_PATH="%(ENV_XRAY_CONFIG_PATH)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[unix_http_server]
file=/var/run/supervisor.sock
; The web process runs as `node` (see [program:web]) but calls
; `supervisorctl restart xray` after key issuance/revocation. Grant the
; node user's group access to the control socket or that call fails with
; permission denied and config changes never take effect until redeploy.
chmod=0770
chown=root:node

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
